# lifeos-bot.service
[Unit]
Description=LifeOS Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw/lifeos
EnvironmentFile=/opt/openclaw.env
Environment="HOME=/home/openclaw"
Environment="PATH=/home/openclaw/lifeos/venv/bin:/home/openclaw/.openclaw/workspace/homebrew/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/openclaw/lifeos/venv/bin/python bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
